# ansible/roles/ojp_proxy/templates/ojp-server.service.j2
# {{ ansible_managed }}
# /etc/systemd/system/ojp-server.service
# OJP Server {{ ojp_version }} — managed by Ansible

[Unit]
Description=OJP Server {{ ojp_version }}
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/java \
    -Duser.timezone=UTC \
    -Dojp.libs.path={{ ojp_install_dir }}/ojp-libs \
    -Dojp.server.port={{ ojp_server_port }} \
    -Dojp.server.virtualThreads.enabled={{ ojp_server_virtual_threads_enabled | lower }} \
    -Dojp.server.slowQuerySegregation.enabled=true \
    -Dojp.server.slowQuerySegregation.classificationMode=ABSOLUTE_THRESHOLD \
    -Dojp.server.slowQuerySegregation.slowQueryThresholdMs=5000 \
    -Dojp.server.slowQuerySegregation.slowSlotTimeout={{ bench_connection_timeout_ms }} \
    -Dojp.server.slowQuerySegregation.fastSlotTimeout={{ bench_connection_timeout_ms }} \
    -Dojp.server.logLevel=WARN \
    -Dojp.telemetry.enabled=false \
    -Dojp.server.circuitBreakerTimeout=5000 \
    -Dojp.server.threadPoolSize=16 \
    -jar {{ ojp_install_dir }}/bin/ojp-server-{{ ojp_version }}-shaded.jar
StandardOutput=append:{{ ojp_log_file }}
StandardError=append:{{ ojp_log_file }}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
